# gala_status_server.service
[Unit]
Description=Gala Node HTTP Status Server
After=network.target

[Service]
ExecStart=/usr/bin/python3 /root/gala-node-setup/status_server.py
WorkingDirectory=/root/gala-node-setup
Restart=always
User=root

[Install]
WantedBy=multi-user.target
